# SmartTooling_ITIS/my_general_startup.service
[Unit]
Description=user pi start general services in different threads
After=network.target

[Service]
WorkingDirectory=/home/pi
Type=idle
User=pi
ExecStart=/usr/bin/python /media/usb/general_startup.py
Restart=always

[Install]
WantedBy=multi-user.target
